# supervisord.conf
[supervisord]
nodaemon=true
user=root
loglevel=info
logfile=/app/logs/supervisor.log

[program:tor]
command=/usr/bin/tor -f /etc/tor/torrc
autostart=true
autorestart=true
user=root
stdout_logfile=/app/logs/tor-stdout.log
stderr_logfile=/app/logs/tor-stderr.log
priority=100
startsecs=20

[program:qnk-node]
command=/app/q-api-server
autostart=true
autorestart=true
user=root
stdout_logfile=/app/logs/qnk-stdout.log
stderr_logfile=/app/logs/qnk-stderr.log
priority=200
startsecs=25
depends_on=tor
environment=RUST_LOG="info,q_api_server=debug,q_network=debug,q_tor_client=debug,q_dns_phantom=debug"