# scripts/phase22-service-beta.service
[Unit]
Description=Q-NarwhalKnight API Server - Phase 22 Full Transition Test & Auto-Purge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/orobit/shared/q-narwhalknight

# ✅ v6.5.0-beta Phase 22 Configuration - Fresh database, genesis Feb 16, 2026
Environment="Q_DB_PATH=./data-mine22"
Environment="Q_NETWORK_ID=testnet-phase22"
Environment="Q_IS_VALIDATOR=true"
Environment="Q_P2P_PORT=9001"
Environment="Q_ENABLE_AI=1"
Environment="RUST_LOG=info"
Environment="Q_ALLOW_SOLO_MINING=true"

# 🚀 v1.3.1-beta: TURBO SYNC with WAL DURABILITY GUARANTEES
Environment="Q_TURBO_SYNC=1"
Environment="Q_BATCHED_WRITES=1"
Environment="Q_STATE_SYNC=1"

# 🛡️ v5.5.4: Pre-flight DB integrity check on every restart (mainnet safety)
Environment="Q_PREFLIGHT_CHECK=1"

# AI Resource Limits (optimized for 18-core CPU)
Environment="Q_AI_THREADS=12"
Environment="Q_AI_MAX_CONCURRENT=2"

# 🔐 v6.5.0-beta: RocksDB Encryption with Phase 22 keys
Environment="Q_ENCRYPTION_KEYS_FILE=/opt/encryption-phase22.keys"
Environment="Q_ENCRYPTION_PASSPHRASE=Qnk-Phase22-TransitionTest-2026-ServerBeta-Prod-Key"

# 🐳 v1.2.2-beta: Docker P2P Sync Optimization - Announce ONLY public IP
Environment="Q_EXTERNAL_ADDRESS=/ip4/185.182.185.227/tcp/9001"
Environment="Q_FILTER_DOCKER_ADDRESSES=true"

# 🌐 v3.4.3-browser: Browser P2P Support - WSS external address for nginx proxy
Environment="Q_EXTERNAL_WSS_ADDRESS=/dns4/quillon.xyz/tcp/9443/wss"

# ⚡ v2.8.0: Ultra-low-latency gossipsub for instant P2P mining reward propagation
Environment="Q_GOSSIPSUB_HEARTBEAT_MS=50"
Environment="Q_GOSSIPSUB_FLOOD_PUBLISH=true"

ExecStart=/opt/orobit/shared/q-narwhalknight/target/release/q-api-server --port 8080
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=q-api-server

# ✅ v0.9.76-beta: Graceful shutdown timeout (5 minutes to prevent SIGKILL data loss)
TimeoutStopSec=300
KillMode=mixed
SendSIGKILL=yes

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
